# yocto/meta-edgewatch/recipes-edgewatch/edgewatch-agent/edgewatch-agent.service
[Unit]
Description=EdgeWatch Agent
After=network-online.target

[Service]
Type=simple
Environment=API_BASE=http://edgewatch.local
ExecStart=/usr/bin/env python3 /usr/bin/edgewatch/agent.py
Restart=always

[Install]
WantedBy=multi-user.target
